# test/fixtures/test_lesson/backend/example/services.conf
[program:redis]
command=redis-server
autorestart=true
